# ops/com.overwhelmreset.operator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.overwhelmreset.operator</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__ROOT__/scripts/run_operator_supervisor.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>__ROOT__</string>
    <key>StandardOutPath</key>
    <string>__ROOT__/.runtime/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__ROOT__/.runtime/launchd.err.log</string>
  </dict>
</plist>
